# gfs2/scripts/82-gfs2-withdraw.rules
SUBSYSTEM=="gfs2", ACTION=="offline", RUN+="/bin/sh /usr/sbin/gfs2_wd_udev.sh"

