# gfs2/scripts/82-gfs2-withdraw.rules.in
SUBSYSTEM=="gfs2", ACTION=="offline", RUN+="/bin/sh @libexecdir@/gfs2_withdraw_helper"

